# debian/postinst
#!/bin/sh
update-rc.d ajenti defaults > /dev/null
#DEBHELPER#

