# debian/postinst
#!/bin/sh
#DEBHELPER#
update-rc.d ajenti defaults > /dev/null
/etc/init.d/ajenti start
