# debian/postrm
#!/bin/sh
update-rc.d -f ajenti remove > /dev/null
#DEBHELPER#

