# debian/postinst
#!/bin/sh
update-rc.d ajenti defaults > /dev/null
/etc/init.d/ajenti start
#DEBHELPER#

