# dist/gentoo/files/.svn/text-base/ajenti.svn-base
#!/sbin/runscript
# Copyright 1999-2010 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	use net
} 

start() {
	ebegin "Starting Ajenti" 
	start-stop-daemon --start --exec /usr/bin/ajenti-panel --pidfile /var/run/ajenti.pid -- --start
	eend $?
}

stop() {
	ebegin "Stopping Ajenti"
	start-stop-daemon --stop --exec /usr/bin/ajenti-panel --pidfile /var/run/ajenti.pid
	eend $?
}
