# packaging/files/ajenti
#!/bin/sh

### BEGIN INIT INFO
# Provides:          ajenti
# Required-Start:    $network $syslog $local_fs
# Required-Stop:     $network $syslog $local_fs
# Should-Start:      $local_fs
# Should-Stop:       $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Ajenti
# Description:       Ajenti administration frontend
### END INIT INFO

. /lib/lsb/init-functions

NAME=Ajenti
DAEMON=/usr/bin/ajenti-panel
PIDFILE=/var/run/ajenti.pid

test -x $DAEMON || exit 5

case "$1" in
	start)
		if [ -e $PIDFILE ]; then
			status_of_proc -p $PIDFILE $DAEMON $NAME && status="0" || status="$?"
			if [ $status = "0" ]; then
				exit
			fi
		fi
		log_daemon_msg "Starting" "$NAME"
		if start-stop-daemon --start --quiet --oknodo --pidfile $PIDFILE --exec $DAEMON -- -d ; then
			log_end_msg 0
		else
			log_end_msg 1
		fi
		;;
	stop)
		if [ -e $PIDFILE ]; then
			status_of_proc -p $PIDFILE $DAEMON "Stoppping $NAME" && status="0" || status="$?"
			if [ "$status" = 0 ]; then
				start-stop-daemon --stop --quiet --oknodo --pidfile $PIDFILE
				/bin/rm -rf $PIDFILE
			fi
		else
		   log_daemon_msg "$NAME is not running"
		   log_end_msg 0
		fi
		;;
	restart)
  		$0 stop && sleep 2 && $0 start
		;;
 	status)
  		if [ -e $PIDFILE ]; then
   			status_of_proc -p $PIDFILE $DAEMON $NAME && exit 0 || exit $?
  		else
   			log_daemon_msg "$NAME is not running"
   			log_end_msg 0
   		fi
  		;;
	*)
		echo "Usage: $0 {start|stop|restart|status}"
		exit 1 
esac

exit 0
